# app-misc/keycodes-hp6510/files/v1.0-etc_init.d_keycodes
#!/sbin/runscript

depend() {
	need localmount
}

start() {
	# Load scan code to key code mappings
	if [ -n "$KEYCODES" -a -x /usr/bin/setkeycodes ]; then
		ebegin "Loading scan code to key code mappings"
		/usr/bin/setkeycodes $KEYCODES
		eend $? "Error loading key code mappings"
	fi
}
